# bundles/001-core/rootfs/etc/profile.d/lang.sh
#!/bin/sh
# Set the system locale.  (no, we don't have a menu for this ;-)
# For a list of locales which are supported by this machine, type:
#   locale -a

#export LANG=en_US.UTF-8
export LANG=tr_TR.UTF-8